# fast_run/service_config/cron/pxw-25/crontab

50 3 * * * rsync -avz --delete --log-file=/var/log/rsync/nas /data/ /501_data/25_data/
